# files/system/etc/profile.d/aos-fastfetch.sh
#!/bin/sh
# Set a.OS default fastfetch alias to use our custom config
alias fastfetch='fastfetch -c /usr/share/a.os/fastfetch.jsonc'